# services/set-random-bg.service
[Unit]
Description=Set a random GNOME desktop background

[Service]
Type=oneshot
ExecStart=/home/eduardotc/bin/set_random_bg.sh
